# linux/w_check.cron
* * * * * echo "" >> ~/w_checks.txt; w -i >> ~/w_checks.txt
